# configurations/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:process_users]
command=/usr/local/bin/middle_earth/processusers
autostart=true
autorestart=true
stderr_logfile=/var/log/process_users.err.log
stdout_logfile=/var/log/process_users.out.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout

